# haas-firewall.timer
[Unit]
Description=Daily Haas Firewall Sync

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
